# deploy/systemd/my-api-drift-check.service
[Unit]
Description=My API Drift Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
ExecStartPre=-/bin/bash -lc '/srv/my-api/current/scripts/self_heal_market_api_proxy.sh'
ExecStart=/bin/bash -lc '/srv/my-api/current/scripts/check_server_drift.sh --alert-on-drift'
