# deploy/systemd/personal-api-admin-backend.service
[Unit]
Description=Personal API Admin Backend
After=network.target
OnFailure=my-api-alert@%p.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/srv/my-api/current/backend
EnvironmentFile=/etc/my-api/backend.env
ExecStart=/srv/my-api/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --proxy-headers
Restart=always
RestartSec=3
TimeoutStopSec=15
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
